# files/system/etc/systemd/system/netxms-server.service.d/override.conf
[Unit]
After=network.target postgres.service